# app/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g 'daemon off;'
autorestart=true

[program:init-letsencrypt]
command=/init-letsencrypt.sh
autorestart=false
